# imgurt.service.d/imgurt.service
[Unit]
Description=Grab random wallpaper

[Service]
Environment="DISPLAY=:0"
ExecStart=/home/evan/resources/imgurt/imgurt.py