# systemd/imgurt.timer
[Unit]
Description=Run random_wallpaper.service periodically

[Timer]
OnCalendar=*-*-* 00:00:0/30
Persistent=true

[Install]
WantedBy=timers.target
